# .cron
*/30 * * * * {cmd} "TukUserTemp::cleanup()"